# initrd/etc/inittab
::sysinit:/etc/rc on
::respawn:-/bin/bash -l
::ctrlaltdel:/etc/rc off
::ctrlaltdel:/busybox/reboot
